# wifi_loggin.supervisor.conf
[program:wifi_login.py]
command=/usr/local/bin/wifi_login.py
autostart=true
autorestart=true
user=root
directory=/tmp

